# systemd.timer
[Unit]
Description=rightmove daily

[Timer]
OnCalendar=*-*-* 11:20:00
Persistent=true

[Install]
WantedBy=timers.target
