# systemd.timer
[Unit]
Description=rightmove daily

[Timer]
OnCalendar=*-*-* 11:20:00
Persistent=true
Unit=rightmove-run.service

[Install]
WantedBy=timers.target
